# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

# Xvfb - Virtual X11 display server
[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1280x1024x24 -ac
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.log

# x11vnc - VNC server for X11 display
[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -forever -nopw -rfbport 5900
environment=DISPLAY=:0
autostart=true
autorestart=true
priority=2
depends_on=xvfb
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.log

# noVNC - Web-based VNC viewer
[program:novnc]
command=/usr/share/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080
autostart=true
autorestart=true
priority=3
depends_on=x11vnc
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.log

# InvestPro Java Application
[program:investpro]
command=java -Xmx2g -Xms512m -Djava.awt.headless=false -jar investpro.jar
environment=DISPLAY=:0,PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
directory=/app
autostart=true
autorestart=true
priority=4
depends_on=xvfb
stdout_logfile=/var/log/supervisor/investpro.log
stderr_logfile=/var/log/supervisor/investpro.log
stopasgroup=true
startsecs=10

# Logs directory
[supervisorctl]

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
